# docker/s6/s6-rc.d/ikaros/run
#!/bin/sh
cd /ikaros
export HOME=/ikaros
exec s6-setuidgid abc flask run --port=12346
