# docker/s6/s6-rc.d/ikaros/run
#!/bin/sh
cd /app
export HOME=/app
exec s6-setuidgid tomoki flask run --port=12346
